# debian_package/DEBIAN/postinst
#!/bin/bash

update-rc.d plsmonitor defaults
# anpassen der Dateiberechtigungen
